# blocks_http_server.service
[Unit]
Description=html server

[Service]
ExecStart=/home/joe/blocks_multiplayer/http_server.py
