# blocks_ws_server.service
[Unit]
Description=websocket server

[Service]
ExecStart=/home/joe/blocks_multiplayer/ws_server.py
